# services/git-daemon/run
#!/command/execlineb -P

/usr/bin/git daemon --reuseaddr --base-path=/repos/ /repos/
